# package/share/solseek/solseek-uc.service
[Unit]
Description=Run solseek update check

[Service]
Type=oneshot
ExecStart=/usr/bin/solseek -cu

[Install]
WantedBy=default.target